# systemd/dipper_ai.service
[Unit]
Description=dipper_ai DDNS manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/etc/dipper_ai
ExecStart=/usr/bin/dipper_ai daemon
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target
